# cli-proxy-api.service
[Unit]
Description=CLI Proxy API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/xxnuo/projects/CLIProxyAPI
ExecStart=/home/xxnuo/projects/CLIProxyAPI/cli-proxy-api
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

